# scripts/systemd/talicska-power.service
[Unit]
Description=Talicska Robot — Power Mode (nvpmodel MAXN_SUPER + jetson_clocks)
After=multi-user.target
Before=talicska-robot.service

[Service]
Type=oneshot
RemainAfterExit=yes
# nvpmodel MAXN_SUPER (mode 2) — a '-' prefix: hiba esetén NEM töri le a service-t
ExecStart=-/usr/sbin/nvpmodel -m 2
# jetson_clocks — maximális órajel minden magon
ExecStart=-/usr/bin/jetson_clocks
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
